# ozon-tracker.service
# ============================================================================
# SYSTEMD СЕРВИС ДЛЯ OZON TRACKER
# ============================================================================
#
# Назначение: Автозапуск приложения при загрузке сервера
#
# Инструкция по установке:
# 1. Замени YOUR_USERNAME на имя пользователя на сервере
# 2. Скопируй: sudo cp ozon-tracker.service /etc/systemd/system/
# 3. Перезагрузи демон: sudo systemctl daemon-reload
# 4. Включи автозапуск: sudo systemctl enable ozon-tracker
# 5. Запусти сервис: sudo systemctl start ozon-tracker
#
# Полезные команды:
# - Проверить статус: sudo systemctl status ozon-tracker
# - Посмотреть логи: sudo journalctl -u ozon-tracker -f
# - Перезапустить: sudo systemctl restart ozon-tracker
# - Остановить: sudo systemctl stop ozon-tracker
#
# ============================================================================

[Unit]
Description=Ozon FBO Tracker - Flask Application
After=network.target

[Service]
Type=exec
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=OZON_APP_DIR
Environment="PATH=OZON_APP_DIR/venv/bin"
ExecStart=OZON_APP_DIR/venv/bin/gunicorn --config gunicorn_config.py ozon_app:app

# Автоматический перезапуск при падении
Restart=always
RestartSec=10

# Graceful shutdown - ждём завершения синхронизации (до 5 минут)
TimeoutStopSec=300
KillMode=mixed
KillSignal=SIGTERM

# Логи
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ozon-tracker

# Безопасность
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
